# deploy/systemd/auto48-worker.service
[Unit]
Description=auto48 arq worker (image processing)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=auto48
Group=auto48
WorkingDirectory=/opt/auto48/app
EnvironmentFile=/opt/auto48/.env
# Phase 1a worker only. The ingest/alerts workers belong to 1b and import
# code that isn't part of this deploy.
ExecStart=/opt/auto48/venv/bin/arq auto48.workers.images.WorkerSettings
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
